# supervisord.omnibus.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:postgres]
command=postgres -D /app/postgres/data -k /var/run/postgresql -h 127.0.0.1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=0

[program:backend]
command=/bin/sh -c "until pg_isready -h 127.0.0.1 -q 2>/dev/null; do sleep 0.5; done && exec /app/backend/.venv/bin/uvicorn main:app --host 127.0.0.1 --port %(ENV_BACKEND_PORT)s"
directory=/app/backend
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=1

[program:frontend]
command=node /app/frontend/dist/server/entry.mjs
directory=/app/frontend
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=2
environment=HOST="0.0.0.0",PORT="7330"
